# credivision-rotation.service
[Unit]
Description=CrediVision Tab Rotation Service
After=credivision-kiosk.service
Wants=credivision-kiosk.service

[Service]
Type=simple
User=informa
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/informa/.Xauthority
ExecStartPre=/bin/sleep 60
ExecStart=/usr/bin/python3 /home/informa/Documentos/CrediVision/rotate_tabs.py 300
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
